# hms-mileage.service
[Unit]
Description=Highland Media Services — Mileage Tracker
After=network.target

[Service]
Type=simple
User=REDACTED
WorkingDirectory=/var/www/html/mileage.REDACTEDmediaservices.com
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hms-mileage
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
